# profile.d/50-openjdk.sh
#! /bin/bash
# -*-sh-*-

# Support a local installation of openjdk

pathprepend /usr/lib/jvm/java-25-openjdk/bin
export JAVA_HOME=/usr/lib/jvm/java-25-openjdk/
