# profile.d/50-sway.sh
#! /bin/bash
# -*-sh-*-
# .profile.d/50-sway.sh
# Time-stamp: 2021-Jul-08 22:29
# Configuration for sway

# Debugging
# set -x

export XDG_CURRENT_DESKTOP=sway


