# profile.d/60-arch.sh
#! /bin/bash
# -*-sh-*-

# Debugging
# set -x

# fzf shell integration
eval "$(fzf --bash)"
