# profile.d/60-bat.sh
# Use bat as the MAN pager
# Time-Stamp: 2023-Jul-13 14:07

export MANPAGER="sh -c 'col -bx | bat -l man -p'"
