# profile.d/60-bat.sh
# Use bat as the MAN pager
# Time-Stamp: 2025-Oct-17 12:00

export MANPAGER="bat --style=plain"
