# profile.d/60-golang.sh
#! /bin/bash
# -*-sh-*-

# Support a local installation of the go language

pathprepend /usr/local/go/bin

