# profile.d/60-rustlang.sh
#! /bin/bash
# -*-sh-*-

# Support a local rust binary path

pathprepend ~/.cargo/bin

