# profile.d/60-sam-ba.sh
#! /bin/bash
# -*-sh-*-

# Support a local installation of the SAM-BA tool

pathprepend /opt/sam-ba_v3.9.1


