# profile.d/60-termswx.sh
#! /bin/bash
#  vim: set ts=4 sw=4 sts=4 tw=120 cc=120 et ft=sh :
# Steen Hegelund
# Time-Stamp: 2026-May-18 10:35
# Adding connection name completion for termswx

source <(termswx --bpaf-complete-style-bash)
