# profile.d/80-nvim.sh
# Support neovim as the default editor
alias vi=nvim
alias vim=nvim
alias vimdiff='nvim -d'
export SYSTEMD_EDITOR=nvim
