# profile.d/80-nvim.sh
# Support neovim as the default editor
alias vim=nvim

