# profile.d/90-docker.sh
# Steen Hegelund
# Time-Stamp: 2023-Jul-10 19:55

pathappend ~/.docker/bin

