# config/systemd/user/i3_wallpaper.timer
[Unit]
Description=Wallpaper Change Timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=2h
Unit=i3_wallpaper.service

[Install]
WantedBy=timers.target
