# config/systemd/user/i3_wallpaper.service
[Unit]
Description=Wallpaper Change service

[Service]
Type=oneshot
ExecStart=%%h/scripts/wallpapers.py -w -z -r %h/Pictures/Wallpapers
